# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:postgresql]
command=su - postgres -c "/usr/lib/postgresql/17/bin/postgres -D /var/lib/postgresql/17/main -c config_file=/etc/postgresql/17/main/postgresql.conf"
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/postgresql.log
stderr_logfile=/var/log/supervisor/postgresql.err

[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis.err

[program:fastapi]
command=/app/backend/.venv/bin/uvicorn edof_netz.main:app --host 127.0.0.1 --port 8000
directory=/app/backend/src
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/supervisor/fastapi.log
stderr_logfile=/var/log/supervisor/fastapi.err

[program:celeryworker]
command=/app/backend/.venv/bin/celery -A edof_netz.tasks worker --loglevel=info
directory=/app/backend/src
autostart=true
autorestart=true
priority=40
stdout_logfile=/var/log/supervisor/celeryworker.log
stderr_logfile=/var/log/supervisor/celeryworker.err

[program:celerybeat]
command=/app/backend/.venv/bin/celery -A edof_netz.tasks beat --loglevel=info --scheduler celery.beat.PersistentScheduler
directory=/app/backend/src
autostart=true
autorestart=true
priority=40
stdout_logfile=/var/log/supervisor/celerybeat.log
stderr_logfile=/var/log/supervisor/celerybeat.err

[program:nextjs]
command=node /app/web/node_modules/.bin/next start --hostname 127.0.0.1 --port 3001
directory=/app/web
autostart=true
autorestart=true
priority=50
stdout_logfile=/var/log/supervisor/nextjs.log
stderr_logfile=/var/log/supervisor/nextjs.err

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=60
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.err
